# crontab
#Run script everyday at 7 00 UTC
0 7 * * * cd /app/ && /usr/local/bin/python3 /app/main.py > /proc/1/fd/1 2>/proc/1/fd/2
